# echo-admin/deployment/supervisor.conf
[program:echo-admin]
command=/opt/echo-admin/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8100
directory=/opt/echo-admin/backend
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/echo-admin.stdout.log
stderr_logfile=/var/log/supervisor/echo-admin.stderr.log
